# hap/installer/linux/comelit-hub-hap.service
[Unit]
Description=Comelit HUB HAP
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/comelit-hub-hap-wrapper.sh
User=comelit
Group=comelit
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/comelit-hub-hap
EnvironmentFile=/etc/comelit-hub-hap/comelit-hub-hap.env

# Allow the service to write PID file (RuntimeDirectory creates /run/comelit-hub-hap)
RuntimeDirectory=comelit-hub-hap
PIDFile=/run/comelit-hub-hap/comelit-hub-hap.pid

# Hardening options
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/comelit-hub-hap /var/log/comelit-hub-hap /var/run

[Install]
WantedBy=multi-user.target
